# cps210/crontab
*/5 * * * * /home/condor/cps210/monitor.sh
